# supervisord.conf
[supervisord]
nodaemon=true
[program:start]
command=bash
directory=/project
autostart=true
startretries=3
user=root
command=/rosa/start.sh